# crontab
0 7-19 * * 0-6 symfony console app:output-products >> log/crontab.log
